# scripts/ops/com.unitares.governance-mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  TEMPLATE: Copy to ~/Library/LaunchAgents/com.unitares.governance-mcp.plist
  and replace all PLACEHOLDER values with your actual paths and credentials.

  Required changes:
    - /PATH/TO/PYTHON3     -> your Python 3.12+ interpreter
    - /PATH/TO/UNITARES    -> your unitares project checkout
    - /YOUR/HOME           -> your home directory
    - ALLOWED_HOSTS/ORIGINS -> your LAN IPs or remove if localhost-only
    - UNITARES_HTTP_API_TOKEN / UNITARES_CONTINUITY_TOKEN_SECRET -> generate with:
        python3 -c "import secrets; print(secrets.token_urlsafe(32))"
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.unitares.governance-mcp</string>

    <key>ProgramArguments</key>
    <array>
        <string>/PATH/TO/PYTHON3</string>
        <string>/PATH/TO/UNITARES/src/mcp_server.py</string>
        <string>--port</string>
        <string>8767</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/PATH/TO/UNITARES</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/PATH/TO/UNITARES/data/logs/mcp_server.log</string>

    <key>StandardErrorPath</key>
    <string>/PATH/TO/UNITARES/data/logs/mcp_server_error.log</string>

    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>4096</integer>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- Resident roster for this deployment (comma-separated labels).
             Each named resident is its own N=1 calibration class. Omit or
             leave empty for a user-agnostic install with no named residents.
             MUST match the residents this deployment actually runs (and the
             UNITARES_RESIDENTS set in each resident's plist). -->
        <key>UNITARES_RESIDENTS</key>
        <string>Lumen,Vigil,Sentinel,Watcher,Steward,Chronicler</string>
        <!-- Resident-progress probe manifest (JSON). Points at the shipped
             canonical fleet; copy + edit for a different deployment, or omit
             for a user-agnostic install that probes no residents. -->
        <key>UNITARES_RESIDENT_PROGRESS_MANIFEST</key>
        <string>/PATH/TO/UNITARES/config/resident_progress.example.json</string>
        <key>UNITARES_BIND_ALL_INTERFACES</key>
        <string>1</string>
        <!-- Host/Origin allowlist (DNS-rebinding protection, on by default).
             A request whose Host header is not listed here is rejected with
             HTTP 403 BEFORE any auth runs. localhost-only is correct for a
             purely local install. To reach the server from a remote MCP
             connector (Claude.ai, Perplexity, etc.) through a reverse proxy
             or Cloudflare tunnel, add the PUBLIC hostname the client connects
             to — e.g. for https://your-tunnel-host.example.org/mcp/ add both
             the bare host and the :* form:
               your-tunnel-host.example.org,your-tunnel-host.example.org:*
             and the matching https origin below. -->
        <key>UNITARES_MCP_ALLOWED_HOSTS</key>
        <string>localhost:*</string>
        <key>UNITARES_MCP_ALLOWED_ORIGINS</key>
        <string>http://localhost:*</string>
        <!-- Static bearer-token gate for /mcp (API-key auth). Empty/unset =
             gate OFF (localhost dev). For a publicly reachable connector set
             this to a long random secret (python3 -c "import secrets;
             print(secrets.token_hex(32))") and paste the SAME value into the
             client's "API key" field. Comma-separated allows rotation. -->
        <key>UNITARES_MCP_BEARER_TOKENS</key>
        <string></string>
        <key>PYTHONPATH</key>
        <string>/PATH/TO/UNITARES</string>
        <key>DB_BACKEND</key>
        <string>postgres</string>
        <key>DB_POSTGRES_URL</key>
        <string>postgresql://postgres:postgres@localhost:5432/governance</string>
        <key>DB_POSTGRES_MIN_CONN</key>
        <string>5</string>
        <key>DB_POSTGRES_MAX_CONN</key>
        <string>25</string>
        <key>UNITARES_KNOWLEDGE_BACKEND</key>
        <string>postgres</string>
        <key>UNITARES_DIALECTIC_BACKEND</key>
        <string>postgres</string>
        <key>HOME</key>
        <string>/YOUR/HOME</string>
        <key>HF_HOME</key>
        <string>/YOUR/HOME/.cache/huggingface</string>
        <key>TRANSFORMERS_OFFLINE</key>
        <string>1</string>
        <key>HF_HUB_OFFLINE</key>
        <string>1</string>
        <key>UNITARES_HTTP_API_TOKEN</key>
        <string>GENERATE_YOUR_OWN_TOKEN</string>
        <key>UNITARES_CONTINUITY_TOKEN_SECRET</key>
        <string>GENERATE_YOUR_OWN_SECRET</string>
        <key>UNITARES_LLM_MODEL</key>
        <string>qwen3.6:27b-coding-nvfp4</string>
    </dict>
</dict>
</plist>
