# scripts/com.unitares.governance-mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.unitares.governance-mcp</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Library/Frameworks/Python.framework/Versions/3.14/bin/python3</string>
        <string>/Users/cirwel/projects/governance-mcp-v1/src/mcp_server.py</string>
        <string>--port</string>
        <string>8767</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/cirwel/projects/governance-mcp-v1</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/cirwel/projects/governance-mcp-v1/data/logs/mcp_server.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/cirwel/projects/governance-mcp-v1/data/logs/mcp_server_error.log</string>
    
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>4096</integer>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/cirwel/projects/governance-mcp-v1</string>
        <key>DB_BACKEND</key>
        <string>postgres</string>
        <key>DB_POSTGRES_URL</key>
        <string>postgresql://postgres:postgres@localhost:5432/governance</string>
        <key>DB_POSTGRES_MIN_CONN</key>
        <string>2</string>
        <key>DB_POSTGRES_MAX_CONN</key>
        <string>10</string>
        <key>UNITARES_KNOWLEDGE_BACKEND</key>
        <string>age</string>
        <key>UNITARES_DIALECTIC_BACKEND</key>
        <string>postgres</string>
        <key>HOME</key>
        <string>/Users/cirwel</string>
        <key>HF_HOME</key>
        <string>/Users/cirwel/.cache/huggingface</string>
        <key>TRANSFORMERS_OFFLINE</key>
        <string>1</string>
        <key>HF_HUB_OFFLINE</key>
        <string>1</string>
    </dict>
</dict>
</plist>
